# deploy/master.service
[Unit]
Description=CF IP Monitor Master
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cfip
WorkingDirectory=/opt/cf_ip_monitor
Environment=PYTHONUNBUFFERED=1
# 依赖通过 `uv sync --frozen --no-dev` 装到 /opt/cf_ip_monitor/.venv
# entrypoint 由 pyproject.toml 的 [project.scripts] 注册为 cf-ip-master
#
# 敏感凭据 (建议放到 /etc/cf-ip-monitor.env, chmod 600, 不入 git)
# config.yaml 里的 ${HUAWEICLOUD_AK} / ${HUAWEICLOUD_SK} / ${HUAWEICLOUD_ZONE_ID} 会自动展开
EnvironmentFile=-/etc/cf-ip-monitor.env
ExecStart=/opt/cf_ip_monitor/.venv/bin/cf-ip-master --config /opt/cf_ip_monitor/config.yaml
Restart=always
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
